# etc/systemd/system/knd-learner.service
[Unit]
Description=Knowdy Learner Service
After=network.target

[Service]
Type=simple
User=knowdy
ExecStart=/usr/bin/knd-learner /etc/knowdy/learner.gsl

[Install]


