# etc/systemd/system/knd_storage.service
[Unit]
Description=Knowdy Storage Service

[Service]
Type=simple
User=knowdy
WorkingDirectory=/home/knowdy
ExecStart=/usr/local/bin/knd_storage /etc/knowdy/storage.ini
Restart=on-failure

[Install]


